# templates/systemd/runrestic.timer.j2
[Unit]
Description=Run runrestic backup

[Timer]
OnCalendar={{ restic_systemd_timer_oncalendar }}
Persistent=true

[Install]
WantedBy=timers.target
